# deploy/solhub-web.service
[Unit]
Description=SOL Hub Web Frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sol-hub/apps/web
ExecStart=/usr/bin/node node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=http://localhost:8000
Environment=NEXT_PUBLIC_S3_URL=http://localhost:9000

[Install]
WantedBy=multi-user.target
